# backend/ai-electro.service
[Unit]
Description=AI Electro Layout Backend
After=network.target

[Service]
User=www-data
WorkingDirectory=%h/ai-electro-layout-starter/backend
Environment="PYTHONUNBUFFERED=1"
ExecStart=%h/ai-electro-layout-starter/backend/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
